# templates/jackett.service.erb
############################################
# WARNING:
# Managed by puppet - DO NOT modify locally.
############################################

[Unit]
Description=Jackett Daemon
After=network.target

[Service]
User=<%= scope['::jackett::user'] %>
Group=<%= scope['::jackett::group'] %>
Type=simple
Restart=always
RestartSec=5
ExecStart=/bin/sh "<%= scope['::jackett::install_path'] %>/jackett_launcher.sh"
TimeoutStopSec=20
#ExecReload=/bin/kill -s HUP $MAINPID

[Install]
WantedBy=multi-user.target

